# system_files/usr/lib/systemd/user/workstation-bootstrap.service
[Unit]
Description=Restore workstation user tools
Documentation=https://github.com/marcortola/workstation-os-image
Wants=network-online.target
Requires=workstation-dms-settings.service
After=network-online.target chezmoi-init.service chezmoi-update.service \
    workstation-dms-settings.service
ConditionPathExists=!%h/.local/state/workstation-os-image/bootstrap-complete

[Service]
Type=oneshot
ExecStart=/usr/libexec/workstation-bootstrap-user

[Install]
WantedBy=graphical-session.target
