# system_files/usr/lib/systemd/user/workstation-dms-settings.service
[Unit]
Description=Initialize captured workstation DMS preference defaults
Documentation=https://github.com/marcortola/workstation-os-image
Requires=dms.service
After=dms.service
Before=workstation-bootstrap.service

[Service]
Type=oneshot
ExecStart=/usr/bin/workstation-apply-dms-settings --initialize

[Install]
WantedBy=graphical-session.target
