# system_files/usr/lib/systemd/user/workstation-claude-mcp-seed.service
[Unit]
Description=Seed default Claude Code MCP servers
Documentation=https://github.com/marcortola/workstation-os-image
Wants=network-online.target
# graphical-session.target: ordered after it for the same reason as
# workstation-bootstrap.service, whose DMS-dependent chain we run after.
After=network-online.target workstation-bootstrap.service graphical-session.target
ConditionPathExists=!%h/.local/state/workstation-os-image/claude-mcp-seeded

[Service]
Type=oneshot
ExecStart=/usr/libexec/workstation-seed-claude-mcp

[Install]
WantedBy=graphical-session.target
