# system_files/usr/lib/systemd/user/workstation-invoice-bookmarks.service
[Unit]
Description=Point the invoices GTK bookmark at the current year/month

[Service]
Type=oneshot
ExecStart=/usr/libexec/workstation-update-invoice-bookmarks
