# system_files/usr/lib/systemd/user/workstation-x11-clipsync.service
[Unit]
Description=Mirror the Xwayland (X11) clipboard into the Wayland clipboard
Documentation=https://github.com/marcortola/workstation-os-image
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/libexec/workstation-x11-clipsync
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
