# system_files/usr/lib/systemd/system/uupd.service.d/10-brew-trust.conf
# Seed linuxbrew's Homebrew tap-trust store before every uupd run so its
# `brew upgrade` (which runs as linuxbrew) does not skip third-party tap
# formulae. workstation-brew-trust.service is Type=oneshot, so uupd's ExecStart
# waits until it has finished. This is the real guarantee (a boot-time
# Before=uupd.service would be inert, since uupd.service is static and only the
# timer pulls it in, including Persistent catch-up runs shortly after boot).
[Unit]
Wants=workstation-brew-trust.service
After=workstation-brew-trust.service
