# system_files/usr/lib/systemd/user/niri.service.d/10-homebrew-path.conf
# Put Homebrew binaries on the graphical session PATH.
#
# niri --session imports its own environment (PATH included) into the systemd
# user manager and the D-Bus activation environment. That import overwrites the
# PATH set by /usr/lib/environment.d, so brew never reaches the session and GUI
# apps (JetBrains Toolbox and the IDEs it launches, etc.) cannot find brew tools
# such as node -- which is why WebStorm reports "Node.js is required" despite
# node being installed. An environment.d drop-in cannot win this race; setting
# niri's own PATH does, because niri then propagates brew to the manager and to
# every app it spawns.
#
# Brew is prepended, mirroring `brew shellenv` (which the interactive fish shell
# already applies via zdots). systemd replaces PATH wholesale, so the non-brew
# entries restate niri's default session PATH; %h expands to the user home so no
# username is hardcoded. Keep the non-brew entries in sync if the base session
# PATH changes upstream.
[Service]
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:%h/.local/share/JetBrains/Toolbox/scripts
