# system_files/usr/lib/systemd/system/workstation-brew-trust.service
[Unit]
Description=Trust third-party Homebrew tap formulae for the linuxbrew user
Documentation=https://github.com/marcortola/workstation-os-image
# Order after brew-setup.service (it untars /home/linuxbrew). Do NOT order after
# brew-proxy-setup.service: it is itself After=multi-user.target while this unit
# is WantedBy=multi-user.target, so ordering after it forms a cycle that systemd
# breaks by dropping this unit's job. The ConditionPathExists on
# proxy/brew-original below instead ensures brew-proxy-setup has already run; a
# boot where it has not yet is a clean no-op that the next boot (or the uupd
# drop-in, which runs this in uupd's own cycle-free transaction) covers.
After=brew-setup.service
# bin/brew must exist (brew provisioned) and proxy/brew-original must exist
# (brew-proxy-setup has swapped the symlink and finished its recursive
# chown/chmod), so this no-ops cleanly until brew is fully provisioned instead
# of racing that setup.
ConditionPathExists=/home/linuxbrew/.linuxbrew/bin/brew
ConditionPathExists=/home/linuxbrew/.linuxbrew/proxy/brew-original

[Service]
Type=oneshot
User=linuxbrew
# brew needs an unconfined SELinux context for its operations, matching
# uupd.service which runs brew the same way.
SELinuxContext=system_u:unconfined_r:unconfined_t:s0
ExecStart=/usr/libexec/workstation-brew-trust

[Install]
WantedBy=multi-user.target
