# system_files/usr/lib/systemd/system/workstation-brew-trust.service
[Unit]
Description=Trust third-party Homebrew tap formulae for the linuxbrew user
Documentation=https://github.com/marcortola/workstation-os-image
# brew-setup.service untars /home/linuxbrew; brew-proxy-setup.service installs
# the bin/brew dispatch wrapper and chowns /home/linuxbrew to linuxbrew. Order
# after both (they are pulled in by their own presets); the ConditionPathExists
# guards below make this a clean no-op if brew is not yet fully provisioned.
After=brew-setup.service brew-proxy-setup.service
# proxy/brew-original exists only after brew-proxy-setup has swapped the symlink
# (and, on first boot, before its recursive chown/chmod). Gating on it means we
# no-op cleanly until brew is fully provisioned rather than racing that setup.
ConditionPathExists=/home/linuxbrew/.linuxbrew/bin/brew
ConditionPathExists=/home/linuxbrew/.linuxbrew/proxy/brew-original

[Service]
Type=oneshot
User=linuxbrew
# brew needs an unconfined SELinux context for its operations, matching
# uupd.service which runs brew the same way.
SELinuxContext=system_u:unconfined_r:unconfined_t:s0
ExecStart=/usr/libexec/workstation-brew-trust

[Install]
WantedBy=multi-user.target
