# system_files/usr/lib/systemd/system/workstation-docker-users.service
[Unit]
Description=Grant interactive workstation users Docker access
Documentation=https://github.com/marcortola/workstation-os-image
After=systemd-sysusers.service
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/workstation-configure-docker-users

[Install]
WantedBy=multi-user.target
